# hunchentoot
#!/bin/sh

# PROVIDE: hunchentoot
# REQUIRE: LOGIN nginx
# KEYWORD: shutdown

. /etc/rc.subr
. /var/lib/hunchentoot/htoot.conf

name="hunchentoot"
rcvar=`set_rcvar`

start_cmd="htoot_exec start"
stop_cmd="htoot_exec stop"
attach_cmd="htoot_exec attach"
swank_cmd="htoot_exec swank"
log_cmd="htoot_exec log"
extra_commands="attach swank log"

htoot_exec ()
{
    su -m $HT_USR -c "$HTOOTD $1"
}

required_files="/var/lib/hunchentoot/htoot.conf"

load_rc_config $name
run_rc_command "$1"
